# manifest/deploy/supervisor/gf2-demo-api.ini
; 生产环境
[program:gf2-demo-api]
environment = 
    GF_GERROR_BRIEF=true,
    GF_GCFG_FILE=config.prod.yaml
directory =/app/gf2-demo
command =/app/gf2-demo/gf2-demo-api
autostart = true
startsecs = 3
user = root
redirect_stderr = true
stdout_logfile =/app/gf2-demo/gf2-demo-api.log

